# sys-files/80-wireless-regdom.rules
SUBSYSTEMS=="usb", SYSFS{idVendor}=="148f", SYSFS{idProduct}=="3070", RUN+="/usr/sbin/iw reg set BO"
